# githubArchive/crontab
* 1 * * * root /data/main.py >> /var/log/cron.log 2>&1
