# githubArchive/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ;
logfile_maxbytes=50MB ;
logfile_backups=10 ;
loglevel=info ;
pidfile=/var/run/supervisord.pid ;
nodaemon=true ; (start in foreground if true;default false)
minfds=1024 ; (min. avail startup file descriptors;default 1024)
minprocs=200 ; (min. avail process descriptors;default 200)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:githubArchive-consumer]
command=/usr/bin/python3 /data/archivesConsume.py
process_name=%(program_name)s
directory=/data
autostart=true
autorestart=true
redirect_stderr=true

[program:simulator-consumer]
process_name=%(program_name)s
command=/usr/bin/python3 /data/simulatorConsume.py
directory=/data
autostart=true
autorestart=true
redirect_stderr=true

[program:elastic-consumer]
command=/usr/bin/python3 /data/elasticSearchConsume.py
process_name=%(program_name)s
directory=/data
autostart=true
autorestart=true
redirect_stderr=true
